# megadrive/light/wls.service
[Unit]
Description=wls
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Environment=HOME=/root
ExecStartPre=-/usr/bin/docker kill wls
ExecStartPre=-/usr/bin/docker rm wls
ExecStartPre=/usr/bin/docker pull faddat/wls-pi
ExecStart=/usr/bin/docker run --name wls -p 5555:6969 -p 8095:8090 -v whaled:/whaled faddat/wls-pi /usr/bin/whaled --data-dir /whaled
ExecStop=/usr/bin/docker stop wls

[Install]
WantedBy=multi-user.target