# megadrive/condenser.service
[Unit]
Description=Condenser

[Service]
ExecStart=/usr/bin/yarn run production
WorkingDirectory=/condenser

[Install]
WantedBy=multi-user.target
